# build/script/risc_debian/saagent.service
[Unit]
Description=RMM-Agent
After=networking.service

[Service]
User=root
Group=root

Type=simple
WorkingDirectory=/usr/local/AgentService
ExecStart=/usr/local/AgentService/cagent
KillMode=process
KillSignal=SIGTERM
Restart=on-abort
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target 
