# build/script/Ubuntu/sawatchdog.conf

description "RMM Watchdog service"

start on runlevel [2345]
stop on runlevel [!2345]

console none

pre-start script
	ls "/usr/local/AgentService/SAWatchdog_Config"
end script

script
	chdir "/usr/local/AgentService/"
	exec  "/usr/local/AgentService/sawatchdog"	
end script




