# build/script/risc_yocto/sawatchdog.service
[Unit]
Description=RMM-Agent

[Service]
User=root
Group=root

Type=sample
WorkingDirectory=/usr/local/AgentService
ExecStart=/usr/local/AgentService/sawatchdog
TimeoutStopSec=20
KillMode=control-group
KillSignal=SIGKILL
Restart=on-abort
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target 
